# files/cron.d/swift-stats-processor
30 * * * * %STACK_USER% /usr/local/bin/swift-log-stats-collector /etc/swift/log-processor.conf 
