# files/cron.d/swift-stats-log-uploader
10 * * * * %STACK_USER% /usr/local/bin/swift-log-uploader /etc/swift/log-processor.conf stats
15 * * * * %STACK_USER% /usr/local/bin/swift-log-uploader /etc/swift/log-processor.conf container-stats
