# files/cron.d/swift-container-stats-log-creator
5 * * * * %STACK_USER% /usr/local/bin/swift-container-stats-logger /etc/swift/log-processor.conf 
